# packages/bterm/profile.d/bterm.sh
alias term='bterm -f /usr/share/bterm/unifont.bgf'
