# packages/backend-service/unit/backend.service
[Unit]
Description=TV streaming server
After=oscam.service network.target
Before=kodi.service

[Service]
ExecStart=/usr/bin/runbackend
Restart=always

[Install]
WantedBy=backend.target
